# netkeeper-interception/init/netcheck.init
#!/bin/sh /etc/rc.common

START=90
NETCHECK_BIN=/usr/bin/netcheck

start() {
	"$NETCHECK_BIN"
}

stop() {
	killall "$NETCHECK_BIN"
}

